# supervisor.conf
[supervisord]
nodaemon=true

[program:http]
command=while true do do_some_simple_job sleep 600 done 

 
